# deploy/ansible/roles/deployment/templates/wracket-session-trim.timer
[Unit]
Description = A timer for Wracket session trimming

[Timer]
OnCalendar = *-*-* 06:00:00
Persistent = true

[Install]
WantedBy = timers.target
