# deploy/ansible/roles/letsencrypt/templates/cert-renew.service
[Unit]
Description = Renew app certificate if necessary

[Service]
Type = oneshot
ExecStart = /bin/bash -c "certbot renew"
User = root
